# configuration/crontab
You can use this one bellow but with no while loop for the script
*/10 * * * * /home/root/monitoring.sh

or this with while loop and sleep in for the script
@reboot sh /home/root/monitoring.sh
