# configuration/crontab
*/10 * * * * /home/username/monitoring.sh
